# packer/open-webui.service
[Unit]
Description=OpenWebUI Service
After=network.target ollama.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/open-webui
ExecStart=/opt/open-webui/venv/bin/open-webui serve
Restart=always

# Basic
Environment=OLLAMA_BASE_URL=http://127.0.0.1:11434
Environment=PORT=8080
Environment=HOST=0.0.0.0

# DATABASE_URL and WEBUI_SECRET_KEY
EnvironmentFile=-/opt/open-webui/.env

[Install]
WantedBy=multi-user.target
